# debian/vyatta-webproxy.postinst.in
#!/bin/bash

prefix=@prefix@
exec_prefix=@exec_prefix@
sysconfdir=@sysconfdir@
bindir=@bindir@
sbindir=@sbindir@
urlfilterdatadir=/opt/vyatta/etc/config/url-filtering

# don't start squid on boot
update-rc.d -f squid remove >/dev/null

# save a copy of default config files
cp /etc/squid/squid.conf /etc/squid/squid.conf.bak
cp /etc/squidguard/squidGuard.conf /etc/squidguard/squidGuard.conf.bak
sed -i 's/nocreate/create 640 proxy adm/' /etc/logrotate.d/squid
sed -i 's/delaycompress//' /etc/logrotate.d/squid
sed -i 's/compress//' /etc/logrotate.d/squid
sed -i 's/etc\/init.d\/squid/etc\/init.d\/squid/' /usr/sbin/update-squidguard
sed -i 's/reload/restart/' /usr/sbin/update-squidguard
chown proxy:adm /var/log/squid

# create local conf placeholder
touch /etc/squid/local.conf

# mv databases to global data directory
if [ ! -d $urlfilterdatadir ]; then
    mkdir -p $urlfilterdatadir
fi
if [ -d /var/lib/squidguard ]; then
    mv /var/lib/squidguard $urlfilterdatadir
fi
if [ -d /var/lib/sitefilter ]; then
    mv /var/lib/sitefilter $urlfilterdatadir
fi
